# configs/chron.conf
[program:chronicler]
command=/usr/bin/python /home/olivneh/projects/glass/chronicler.py
